# docker/etc/kspeeder.conf
[program:kspeeder]
command=/bin/sh -c "source /data/local/docker/docker.env && exec /data/local/docker/bin/kspeeder -cachePath ${DISK_CACHE}/Kspeeder"
autostart=true
autorestart=true
