# docker/etc/kspeeder.conf
[program:kspeeder]
command=/data/local/docker/bin/kspeeder \
    -cachePath /mnt/media_rw/048a3593-e231-438c-a980-38149429f948/dockercache
