# docker/etc/kspeeder.conf
[program:kspeeder]
command=/data/local/docker/bin/kspeeder \
    -cachePath %(ENV_DISK_CACHE)s/Kspeeder
